# server/service_files/bulk_data_provider.service
[Unit]
Description=Bulk data provider service for apartment climate monitor
After=network.target

[Service]
Type=simple
User=__USERNAME__
EnvironmentFile=/home/__USERNAME__/climate_monitor/server/.env
WorkingDirectory=/home/__USERNAME__/climate_monitor/server
ExecStart=/home/__USERNAME__/.local/bin/poetry run gunicorn bulk_data_provider:app --bind 127.0.0.1:8000
Restart=always

[Install]
WantedBy=multi-user.target
