# server/service_files/mqtt_postgres_bridge.service
[Unit]
Description=MQTT to Postgres Bridge (Poetry)
After=network.target postgresql.service mosquitto.service
Wants=postgresql.service mosquitto.service
ConditionPathExists=__SERVER_DIR__/.env

[Service]
Type=simple
User=__USERNAME__
EnvironmentFile=__SERVER_DIR__/.env
WorkingDirectory=__SERVER_DIR__
ExecStart=__SERVER_DIR__/.venv/bin/python -u mqtt_postgres_bridge.py
Restart=always
RestartSec=5
TimeoutStartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
