# server/service_files/mqtt_postgres_bridge.service
[Unit]
Description=MQTT to Postgres Bridge (Poetry)
After=network.target

[Service]
Type=simple
User=__USERNAME__
EnvironmentFile=/home/__USERNAME__/climate_monitor/server/.env
WorkingDirectory=/home/__USERNAME__/climate_monitor/server
ExecStart=/home/__USERNAME__/.local/bin/poetry run python mqtt_postgres_bridge.py
Restart=always

[Install]
WantedBy=multi-user.target