# server/service_files/summarizer.service
[Unit]
Description=Prepare data for webpage
After=network.target postgresql.service
Wants=postgresql.service
ConditionPathExists=__SERVER_DIR__/.env

[Service]
Type=oneshot
User=__USERNAME__
EnvironmentFile=__SERVER_DIR__/.env
WorkingDirectory=__SERVER_DIR__
ExecStart=__SERVER_DIR__/.venv/bin/python summarizer.py
StandardOutput=journal
StandardError=journal
